# crontab
@reboot python /initial_run.py
0 4 * * * python /strategy_updater.py